# parking-monitor.service
[Unit]
Description=Smart Parking Monitor - AI-powered parking space detection
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# IMPORTANT: If you chose a username other than 'pi' during Pi Imager setup,
# replace every occurrence of 'pi' below with your actual username.
# Run: whoami   on the Pi to confirm your username.
User=pi
Group=pi
WorkingDirectory=/home/pi/smart-parking-monitor
Environment=PATH=/home/pi/smart-parking-monitor/venv/bin:/usr/bin:/bin
EnvironmentFile=/home/pi/smart-parking-monitor/.env
ExecStart=/home/pi/smart-parking-monitor/venv/bin/python main.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=parking-monitor

[Install]
WantedBy=multi-user.target
