# parking-monitor.service
[Unit]
Description=Smart Parking Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/parking_monitor
ExecStart=/home/pi/parking_monitor/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=parking-monitor

[Install]
WantedBy=multi-user.target
